# configs/supervisor_daphne.conf
[program:daphne]
command = /webapps/chat/env/bin/daphne_start
user = chat
stdout_logfile = /webapps/chat/logs/daphne.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8