# configs/supervisor_chat.conf
[program:chat]
command = /webapps/chat/env/bin/gunicorn_start
user = chat
stdout_logfile = /webapps/chat/logs/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8